# packages/os/linux/tails/config/chroot_local-includes/etc/systemd/user/elizaos-renderer.service
[Unit]
Description=elizaOS local app renderer
ConditionUser=1000
ConditionPathExists=!/run/elizaos/persistence-maintenance
After=elizaos-agent.service
Wants=elizaos-agent.service
RefuseManualStop=yes
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/usr/local/lib/elizaos/start-elizaos-renderer-user
Restart=always
RestartSec=2
SyslogIdentifier=elizaos-renderer
KillMode=mixed
Nice=5
NoNewPrivileges=yes

[Install]
WantedBy=default.target
