# packages/benchmarks/qwen-claw-bench/data/qwenclawbench-v1.1-100/assets/task_00034_server_integration_check_skill_for_daily_stock_analysis/server/services/market-scraper.service
[Unit]
Description=Market Sentiment Scraper
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=admin
Group=docker
WorkingDirectory=/home/admin/server/deploy
ExecStart=/usr/bin/docker compose up market_scraper
ExecStop=/usr/bin/docker compose stop market_scraper
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
