# packages/os/linux/tails/config/chroot_local-includes/usr/lib/systemd/system/user@1000.service.d/oom-memory-pressure-limit.conf
[Service]
# The default limit of 50% caused systemd-oomd to frequently kill
# gnome-shell in our automated test suite.
ManagedOOMMemoryPressureLimit=65%
