# packages/os/linux/elizaos/config/includes.chroot/etc/systemd/system/serial-getty@ttyS0.service.d/autologin.conf
# elizaOS Linux — serial-console autologin (riscv64/arm64 boot with
# console=ttyS0/ttyAMA0). Same rationale as the tty1 override: agetty
# --autologin runs `login -f user`, which needs the live-config-created
# `user` account to already exist, so order this getty After live-config to
# avoid an "Authentication failure" on the serial line.
[Unit]
After=live-config.service

[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin user --keep-baud 115200,38400,9600 --noclear %I $TERM
