# packages/os/linux/tails/config/chroot_local-includes/usr/lib/systemd/system/gnome-session.target
[Unit]
Description=System target for GNOME Session
