# packages/os/linux/tails/config/chroot_local-includes/usr/lib/systemd/system/run-nosymfollow.mount
# This mount can be used to prevent symlink attacks when accessing files.
# Refs: https://gitlab.tails.boum.org/tails/tails/-/merge_requests/339#note_169215

[Unit]
Description=Bind mount of the filesystem root mounted with nosymfollow

[Mount]
What=/
Where=/run/nosymfollow
ReadWriteOnly=yes
DirectoryMode=700
Options=bind,nosymfollow
