# packages/os/linux/tails/config/chroot_local-includes/usr/lib/systemd/system/tails-detect-disk-ioerrors.service
[Unit]
Description=Detect Squashfs and I/O errors
After=systemd-journald.target
Before=gdm.service

[Service]
Type=notify
ExecStart=/usr/local/sbin/tails-detect-disk-ioerrors

[Install]
WantedBy=multi-user.target
