# packages/os/linux/tails/config/chroot_local-includes/usr/lib/systemd/user/tails-additional-software-install.service
# This user unit triggers installation of additional software packages after
# the desktop has started by starting the identically named system unit.

[Unit]
Description=Trigger installation of Additional Software Packages
Documentation=file:///usr/share/doc/elizaos/website/doc.en.html

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/systemctl --system start tails-additional-software-install.service
TimeoutStartSec=0

[Install]
WantedBy=desktop.target
