# packages/os/linux/tails/config/chroot_local-includes/usr/lib/systemd/user/tails-post-greeter-whisperback.service
[Unit]
Description=Report errors detected at the Welcome Screen
Documentation=file:///usr/share/doc/elizaos/website/doc.en.html
ConditionUser=1000
ConditionPathExists=/var/lib/gdm3/post-greeter-whisperback.json

[Service]
Type=oneshot
ExecStart=whisperback /var/lib/gdm3/post-greeter-whisperback.json
RemainAfterExit=yes

[Install]
WantedBy=desktop.target
