# packages/os/linux/tails/config/chroot_local-includes/usr/lib/systemd/user/tails-create-persistent-storage.service
[Unit]
Description=Create a Persistent Storage (if the user chose to do so in the Welcome Screen)
ConditionUser=1000
ConditionPathExists=/var/lib/live/config/tails.create-persistence

[Service]
Type=oneshot
Restart=no

StandardError=journal
ExecStart=/usr/local/lib/elizaos/create-persistent-storage-session

[Install]
WantedBy=desktop.target
