# packages/os/linux/tails/config/chroot_local-includes/usr/lib/systemd/user/tails-report-disk-ioerrors.service
[Unit]
Description=Report errors about disk being corrupted
ConditionUser=1000

[Service]
Type=oneshot
ExecStart=/usr/local/lib/tails-report-disk-ioerrors
RemainAfterExit=yes
