# packages/os/linux/tails/config/chroot_local-includes/usr/lib/systemd/user/tails-report-mac-spoofing-failed.service
[Unit]
Description=Report MAC address spoofing failure
Documentation=file:///usr/share/doc/elizaos/website/doc.en.html
ConditionUser=1000
ConditionPathExists=/var/lib/live/config/tails.mac-spoofing-failed

[Service]
ExecStart=/usr/local/lib/tails-report-mac-spoofing-failed
RemainAfterExit=yes

[Install]
WantedBy=desktop.target
