# packages/os/linux/tails/config/chroot_local-includes/usr/lib/systemd/user/tails-report-disk-partitioning-errors.service
[Unit]
Description=Report possible errors about disk partitioning
Documentation=file:///usr/share/doc/elizaos/website/doc.en.html
ConditionUser=1000
ConditionPathExists=/var/lib/live/config/tails.disk-partitioning-errors

[Service]
ExecStart=/usr/local/lib/tails-report-disk-partitioning-errors
RemainAfterExit=yes

[Install]
WantedBy=desktop.target
