# packages/os/linux/tails/config/chroot_local-includes/usr/lib/systemd/user/tails-report-disk-ioerrors.path
[Unit]
Description=Report errors about disk being corrupted
ConditionUser=1000

[Path]
PathModified=/var/lib/live/tails.disk.ioerrors

[Install]
WantedBy=desktop.target
