# packages/os/linux/tails/config/chroot_local-includes/usr/lib/systemd/user/dbus.service.d/managed-oom-preference.conf
[Service]
# We've seen systemd-oomd kill dbus.service, which leaves the system in
# an unusable state (tails/tails#20631).
ManagedOOMPreference=avoid
