# packages/os/linux/elizaos/config/includes.chroot/etc/systemd/system/getty@tty1.service.d/autologin.conf
# elizaOS Linux — console autologin fallback for headless/no-GDM boots
# (e.g. serial-console debugging on riscv64). agetty --autologin runs
# `login -f user`, which skips PAM password auth but still requires the
# account to exist. live-config creates the `user` account at boot, so this
# drop-in orders the getty After live-config to avoid the "Authentication
# failure" that occurs when `login -f user` runs before the account exists.
# Username matches GDM AutomaticLogin (etc/gdm3/daemon.conf).
#
# On a normal desktop boot GDM owns the session via graphical.target, so this
# getty stays a fallback and is not the path the user lands in.
[Unit]
After=live-config.service

[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin user --noclear %I $TERM
