# packages/os/linux/tails/config/chroot_local-includes/usr/lib/systemd/system/ferm.service.d/no-drop-on-stop.conf
# Otherwise there is a window of time when shutting down when there is
# no firewall (tails#20536).

[Service]
ExecStop=
