# packages/deploy/systemd/units/eliza-probe.service
[Unit]
Description=Eliza — active health probe
Documentation=https://github.com/elizaOS/eliza

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/eliza/env
ExecStart=%h/bin/eliza-health-probe.sh
SyslogIdentifier=eliza-probe
