# packages/os/linux/tails/config/chroot_local-includes/etc/systemd/system/elizaos-root-mode.service
[Unit]
Description=Restore traversable elizaOS live root permissions
DefaultDependencies=no
Requires=run-nosymfollow.mount
After=run-nosymfollow.mount systemd-tmpfiles-setup.service
Before=sysinit.target basic.target dbus.service polkit.service gdm.service tails-persistent-storage.service
ConditionPathIsMountPoint=/run/nosymfollow

[Service]
Type=oneshot
ExecStart=/bin/chmod 0755 / /run/nosymfollow
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
