# packages/os/linux/tails/config/chroot_local-includes/etc/systemd/system/elizaos.path
[Unit]
Description=Start elizaOS when the live user session is ready
RefuseManualStop=yes

[Path]
PathExists=/run/user/1000/bus
Unit=elizaos.service

[Install]
WantedBy=multi-user.target
