# packages/os/linux/tails/config/chroot_local-includes/usr/lib/systemd/user/org.gnome.Shell@wayland.service.d/managed-oom-preference.conf
[Service]
# We've seen systemd-oomd kill gnome-shell, which leaves the system in
# an unusable state.
ManagedOOMPreference=avoid
