# packages/os/linux/elizaos/config/includes.chroot/etc/systemd/system/elizaos-first-boot.service
[Unit]
Description=elizaOS first-boot bootstrap
After=local-fs.target systemd-tmpfiles-setup.service
Before=multi-user.target
ConditionPathExists=!/var/lib/elizaos/.first-boot-complete

[Service]
Type=oneshot
ExecStart=/usr/local/lib/elizaos/first-boot.sh
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
